# src/php/build/modules/php/files/etc/profile.d/phpcs.sh
#!/usr/bin/env bash

export PATH="/usr/local/src/PHP_CodeSniffer/scripts:${PATH}"
