# src/php/build/modules/php/files/etc/profile.d/phpfarm.sh
#!/usr/bin/env bash

export PATH="/usr/local/bin/phpfarm/inst/bin:/usr/local/bin/phpfarm/inst/current/bin:/usr/local/bin/phpfarm/inst/current/sbin:${PATH}"
