# src/php-5.2/run/modules/run/templates/memcached.conf.erb
[program:memcached]
command = /usr/bin/socat TCP-LISTEN:11211,fork TCP:<%= @memcached_host %>:<%= @memcached_port %>
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
