# build/etc/profile.d/phpfarm.sh
#!/usr/bin/env bash

export PATH=/phpfarm/inst/bin:/phpfarm/inst/current-bin:$PATH
