# build/etc/profile.d/phpfarm.sh
#!/usr/bin/env bash

export PATH=/opt/phpfarm/inst/bin:/opt/phpfarm/inst/current-bin:$PATH
