# src/build/modules/php/files/etc/supervisor/conf.d/php.conf
[program:php]
command=/phpfarm/inst/php-5.2.17/sbin/php-fpm start
autostart=true
autorestart=true
