# src/build/modules/php/files/etc/profile.d/phpfarm.sh
#!/usr/bin/env bash

export PATH="/phpfarm/inst/bin:/phpfarm/inst/current/bin:/phpfarm/inst/current/sbin:${PATH}"
