# src/php/build/modules/build/files/etc/profile.d/phpfarm.sh
#!/usr/bin/env bash

export PATH="/usr/local/src/phpfarm/inst/bin:/usr/local/src/phpfarm/inst/current/bin:/usr/local/src/phpfarm/inst/current/sbin:${PATH}"
