# build/etc/supervisor/conf.d/php-5.5.15.conf
[program:php-5.5.15]
command=/opt/phpfarm/inst/php-5.5.15/sbin/php-fpm -F
autostart=true
autorestart=true