# build/etc/supervisor/conf.d/php-5.2.17.conf
[program:php-5.2.17]
command=/opt/phpfarm/inst/php-5.2.17/sbin/php-fpm start