# build/etc/supervisor/conf.d/php-5.3.28.conf
[program:php-5.3.28]
command=/opt/phpfarm/inst/php-5.3.28/sbin/php-fpm -F
autostart=true
autorestart=true