# src/run/modules/php/templates/memcached.conf.erb
[program:memcached]
command=socat TCP-LISTEN:11211,fork TCP:<%= @memcached_host %>:<%= @memcached_port %>
autostart=true
autorestart=true
