# src/run/modules/php/templates/mysqld.conf.erb
[program:mysqld]
command = /usr/bin/socat UNIX-LISTEN:/var/run/mysqld/mysqld.sock,fork,reuseaddr,unlink-early,user=www-data,group=www-data,mode=777 TCP:<%= @db_host %>:<%= @db_port %>
autostart = true
autorestart = true
