# src/run/modules/php/templates/mysqld.conf.erb
[program:mysqld-1]
command = /usr/bin/socat UNIX-LISTEN:/var/run/mysqld/mysqld.sock,fork,reuseaddr,unlink-early,user=www-data,group=www-data,mode=777 TCP:<%= @mysqld_host %>:<%= @mysqld_port %>
autostart = true
autorestart = true

[program:mysqld-2]
command = /usr/bin/socat TCP-LISTEN:<%= @mysqld_port %>,fork TCP:<%= @mysqld_host %>:<%= @mysqld_port %>
autostart = true
autorestart = true
