# src/run/modules/php/templates/redis.conf.erb
[program:redis]
command=socat TCP-LISTEN:6379,fork TCP:<%= @redis_host %>:<%= @redis_port %>
autostart=true
autorestart=true
