# src/php/build/modules/php/files/etc/supervisor/conf.d/postfix.conf
[program:postfix]
command = /etc/init.d/postfix start
startsecs = 0
autostart = true
autorestart = false
stdout_events_enabled = true
stderr_events_enabled = true
