# src/php-5.2/build/modules/build/files/etc/supervisor/conf.d/logs_cron.conf
[program:logs_cron]
command = /usr/bin/tail -F /var/log/cron.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
