# src/php-5.2/build/modules/build/files/etc/supervisor/conf.d/php.conf
[program:php]
command = /usr/bin/pidproxy /usr/local/src/phpfarm/inst/php-5.2.17/logs/php-fpm.pid /usr/local/src/phpfarm/inst/php-5.2.17/sbin/php-fpm start
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
