# src/php-5.2/run/modules/run/files/etc/supervisor/conf.d/crontab_generator.conf
[program:crontab_generator]
command = /usr/local/bin/crontab_generator
autostart = true
startsecs = 0
